# autostart/btsync.desktop
#!/usr/bin/env xdg-open
[Desktop Entry]
Exec=/usr/local/bin/btsync
Icon=system-run
Path=
Terminal=false
Type=Application
