# autostart/ssh-add.desktop
#!/usr/bin/env xdg-open
[Desktop Entry]
Exec=ssh-add < /dev/null
Icon=system-run
Path=
Terminal=false
Type=Application
